# sensor-setup-stage/03-hotspot/files/hostapd-stopper.service
[Unit]
Description=Disable hostapd after 30 minutes, if it was started

[Service]
ExecStart=/usr/bin/orfs-ap-off.sh